# build/files/supervisord.conf
; Sample supervisor config file.

[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=debug                ; (log level;default info; others: debug,warn,trace)


[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface



[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock 

[program:zerotier]
command=zerotier-one /zt-data
autostart=true
redirect_stdout=true
redirect_stderr=true

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
redirect_stdout=true
redirect_stderr=true

[program:testservice]
command=/opt/testservice.sh
autostart=false
redirect_stdout=true
redirect_stderr=true

[program:standalone]
command=/opt/standalone.sh
autostart=false
redirect_stdout=true
redirect_stderr=true
